# build/scripts/hastur-retrieval.conf
# Hastur retrieval Upstart script
#
# Hastur retrieval service for writing incoming data to Cassandra.

description "Hastur Retrieval"

start on [2345]
stop on runlevel[!2345]

respawn
respawn limit 10 5

script
  exec start-stop-daemon --chdir "/opt/hastur-server" --start --exec "/opt/hastur-server/start_retrieval.sh"
end script
