# build/scripts/hastur-core.conf
# Hastur core Upstart script
#
# Hastur core service for writing incoming data to Cassandra.

description "Hastur Core"

start on [2345]
stop on runlevel[!2345]

respawn
respawn limit 10 5
expect fork

exec start-stop-daemon --chdir "/opt/hastur-server" --start --chuid role-hastur:role-hastur --exec "/opt/hastur-server/start_core.sh"
